# config/supervisord.conf
[supervisord]
nodaemon=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
loglevel=DEBUG

[program:rails]
directory=/app/backend/
command=bundle exec rails s -p 3000 -b 0.0.0.0
numprocs=1
autostart=true
autorestart=true

[program:sidekiq]
directory=/app/backend/
command=bundle exec sidekiq -q mailers
numprocs=1
autostart=true
autorestart=true

[program:postfix]
directory=/app/
command=service postfix start
numprocs=1
autostart=true